# systemd/zeicoin-error-monitor.service
[Unit]
Description=ZeiCoin Error Monitor
After=postgresql.service zeicoin-mining.service
Requires=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/zeicoin
EnvironmentFile=-/root/zeicoin/.env.testnet
EnvironmentFile=-/root/zeicoin/.env.local
Environment="ZEICOIN_DB_HOST=127.0.0.1"
Environment="ZEICOIN_DB_NAME=zeicoin_testnet"
Environment="ZEICOIN_DB_USER=zeicoin"
Environment="ZEICOIN_MONITOR_MIN_SEVERITY=MEDIUM"
Environment="ZEICOIN_MONITOR_SERVICE=zeicoin-mining.service"

ExecStart=/root/zeicoin/zig-out/bin/zeicoin_error_monitor
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
